# code/modules/blackbox_exporter/templates/blackbox_exporter.service.erb
[Unit]
Description=Black Box Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/blackbox_exporter --web.listen-address=<%= @bb_listening_port %>

[Install]
WantedBy=multi-user.target
